# build/ipset-blacklists-service_2.0.4_all/DEBIAN/postinst
#!/bin/bash
set -e

# Reload systemd
systemctl daemon-reload

# Create ipset if not exists
if ! ipset list blacklist-inbound >/dev/null 2>&1; then
    ipset create blacklist-inbound hash:net
fi

echo "Installation complete. Start with:"
echo "  sudo systemctl start ipset-blacklist-boot.service"
echo "  sudo systemctl enable --now ipset-blacklist-update.timer"
