# systemd/system/ipset-blacklist-update.timer
[Unit]
Description=ipset Blacklist Update Timer
Documentation=https://github.com/somnium78/ipset-blacklists-service
Requires=ipset-blacklist-update.service

[Timer]
# Run every 4 hours
OnCalendar=*-*-* 00,04,08,12,16,20:00:00
# Run 5 minutes after boot if missed
OnBootSec=5min
# Randomize start time by up to 10 minutes to avoid server load spikes
RandomizedDelaySec=600
# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
